# debian/prerm
#!/bin/sh
# Disable services

systemctl disable gecos-first-login.service
