# deb_dist/gecos-agent-1.0/debian/gecos-agent/etc/init/gecos-first-login.conf
start on desktop-session-start
script
    __firstlogin_prefix__='/usr/local'
    cd "${__firstlogin_prefix__}/bin"
    ./firstlogin-dbusservice 2>/tmp/gecos-first-login-dbus-service.log
end script
